# systemd-service.service
[Unit]
Description=Akane Production Environment
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/cubilizer/AkaneDoThis
ExecStart=/home/cubilizer/AkaneDoThis/startup.sh
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target